# koresh.service
[Unit]
Description=Just Your Koresh
After=syslog.target
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/wackloner/bot-koresh
ExecStart=python3 /home/wackloner/bot-koresh/main.py

RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
